# auto_git.service
[Unit]
Description=Auto Git Commit and Push Service
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/your-repo-path
ExecStart=/usr/bin/python3 /home/%i/your-repo-path/auto_git.py --interval 30
Restart=always
RestartSec=10

# Environment variables for git
Environment="HOME=/home/%i"
Environment="USER=%i"

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target